# dockerfiles/openstack/juno/glance/conf/supervisor-glance.conf
[program:glance-registry]
user=glance
command=/usr/bin/glance-registry
startretries=5
startsecs=5

[program:glance-api]
user=glance
command=/usr/bin/glance-api
startretries=5
startsecs=5
